# cmd-api/apexcore-autodeploy.timer
[Unit]
Description=ApexCore auto-deploy every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=apexcore-autodeploy.service

[Install]
WantedBy=timers.target
