# cmd-api/apexcore-poller.timer
[Unit]
Description=ApexCore Deploy Poller — every 30 minutes

[Timer]
OnBootSec=60s
OnUnitActiveSec=30min
AccuracySec=60s

[Install]
WantedBy=timers.target
